# 17/borgbackup.timer
[Unit]
Description=Borg backup timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=borgbackup.service

[Install]
WantedBy=multi-user.target
